# templates/cron.d/backup.erb
0 0 * * 0           root /opt/backup/backup full
0 0 * * 1-6         root /opt/backup/backup
